# supervisord.conf
[supervisord]
nodaemon=true
stdout_logfile=/dev/stdout
redirect_stderr=true
logfile_maxbytes=0

[program:sserver]
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=100
environment = SS_PORT="7777",KCP_PORT="8888"
command=ss-server -s 0.0.0.0 -p %(ENV_SS_PORT)s -k 299792458 -m rc4-md5 --fast-open -u

[program:kcptun]
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=200
environment = SS_PORT="7777",KCP_PORT="8888"
command=server_linux_amd64 -t "127.0.0.1:%(ENV_SS_PORT)s" -l "0.0.0.0:%(ENV_KCP_PORT)s" -crypt salsa20 -key 299792458 -sndwnd 2048 -rcvwnd 2048 -mode manual -nodelay 1 -resend 2 -nc 1 -interval 20 -nocomp -dscp 46  -mtu 1400 -datashard 70 -parityshard 30
